# scripts/hotel-worker.service
[Unit]
Description=Hotel Rate Sentinel Worker
After=network.target

[Service]
User=successofmentors
Group=successofmentors
WorkingDirectory=/home/successofmentors/.gemini/antigravity/scratch/hotel
ExecStart=/home/successofmentors/.gemini/antigravity/scratch/hotel/venv/bin/celery -A backend.celery_app worker --loglevel=info -Q celery --pool=solo -n worker_turbo@%h
Restart=always
Environment=PYTHONPATH=/home/successofmentors/.gemini/antigravity/scratch/hotel
EnvironmentFile=/home/successofmentors/.gemini/antigravity/scratch/hotel/.env

[Install]
WantedBy=multi-user.target
